# deploy/polycopier.service
[Unit]
Description=Polycopier — Automated Polymarket Copy Trading Bot
Documentation=https://github.com/cbaezp/polycopier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as a dedicated low-privilege user (create with: useradd -r -s /bin/false polycopier)
User=polycopier
Group=polycopier

# All paths relative to the working directory.
# The .env file must live here.
WorkingDirectory=/opt/polycopier

# --headless: skip the TUI, log to stdout (captured by journald)
ExecStart=/opt/polycopier/polycopier --headless

# Restart automatically on crash or OOM kill.
Restart=on-failure
RestartSec=10s

# Give the bot 30s to finish any in-flight orders on SIGTERM.
TimeoutStopSec=30

# Send logs to journald (view with: journalctl -u polycopier -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polycopier

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
# The bot needs to write .env when settings change via [s] at restart
ReadWritePaths=/opt/polycopier
PrivateTmp=true

[Install]
WantedBy=multi-user.target
